# mkosi.extra/etc/systemd/system/nvidia-fabricmanager.service.d/override.conf
[Unit]
# Only start fabricmanager if NVSwitch is present (multi-GPU HGX systems)
ConditionPathExists=/dev/nvidia-nvswitch0
