# mkosi.extra/etc/systemd/system/nvidia-fabricmanager.service.d/override.conf
[Unit]
# Only start fabricmanager on multi-GPU systems (nvidia1 = 2nd GPU exists)
ConditionPathExists=/dev/nvidia1