# mkosi.extra/etc/systemd/system/tinfoil-vllm.service
[Unit]
Description=vLLM
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/tinfoil-vllm

[Install]
WantedBy=multi-user.target
