# mkosi.extra/etc/systemd/system/tinfoil-boot.service
[Unit]
Description=Tinfoil Boot Service
After=network-online.target docker.service containerd.service nvidia-persistenced.service nvidia-fabricmanager.service
Requires=docker.service containerd.service
Wants=nvidia-persistenced.service nvidia-fabricmanager.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/tinfoil-boot
ExecStop=/usr/local/bin/tinfoil-boot shutdown
User=root
Group=root
StandardOutput=journal+console
StandardError=journal+console
# Give plenty of time for cleanup during shutdown
TimeoutStopSec=600

[Install]
WantedBy=default.target
