# mkosi.extra/etc/systemd/system/containerd.service.d/override.conf
[Unit]
After=tinfoil-ramdisk.service nvidia-persistenced.service nvidia-fabricmanager.service
Requires=tinfoil-ramdisk.service
Wants=nvidia-persistenced.service nvidia-fabricmanager.service

