# mkosi.extra/etc/systemd/system/nvidia-cdi-refresh.service.d/override.conf
[Unit]
After=nvidia-persistenced.service nvidia-fabricmanager.service
Wants=nvidia-persistenced.service nvidia-fabricmanager.service